# systemd/pull-axis-recordings.timer
[Unit]
Description=Run Axis recording pull every morning

[Timer]
OnCalendar=*-*-* 08:15:00
Persistent=true

[Install]
WantedBy=timers.target
