# systemd/pull-axis-recordings.service
[Unit]
Description=Pull completed Axis recordings to SSD and generate heartbeat
Wants=network-online.target home-admin-fieldcam-state.mount home-admin-fieldcam-logs.mount
After=network-online.target home-admin-fieldcam-state.mount home-admin-fieldcam-logs.mount

[Service]
Type=oneshot
User=admin
Group=admin
WorkingDirectory=/home/admin/fieldcam
EnvironmentFile=/home/admin/fieldcam/config/heartbeat.env
ExecStart=/home/admin/fieldcam/scripts/run_pull_then_heartbeat.sh
TimeoutStartSec=21600
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
